# logger/simulo-logger.timer
[Unit]
Description=Simulo Log Service Timer

[Timer]
OnCalendar=*:0/20
Persistent=true

[Install]
WantedBy=timers.target
