# logger/simulo-logger.timer
[Unit]
Description=Simulo Log Service Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
